# app/home/config/systemd/user/stop-night-shift-timer.service
[Unit]
Description=Stops gnome-night-shift timer on failure

[Service]
Type=oneshot
ExecStart=/usr/bin/systemctl --user stop gnome-night-shift.timer
